# files/etc/systemd/system/multi-user.target.wants/provision.service
[Unit]
Description=provision
After=network.target auditd.service

[Service]
ExecStart=/bin/provision.sh
Type=Simple

[Install]
WantedBy=multi-user.target
