# bootstrap/etc/systemd/system/multi-user.target.wants/startup.service
[Unit]
Description=startup
After=network.target auditd.service

[Service]
ExecStart=/bin/startup.sh
Type=Simple

[Install]
WantedBy=multi-user.target
